# link.sh
#!/bin/sh

#chkconfig: - 85 15
#description: my game server

ln -s /root/gdmj/server /root/mjserver/game-server/app/servers/pkroom/games/gdmj
ln -s /root/gdmj/web    /root/mjserver/web-server/public/gdmj



